# supervisord.conf
[supervisord]
nodaemon=false

[program:main]
command=./main.py
autorestart=true