# supervisord.conf
[supervisord]
nodaemon=false

[program:proxy]
command=./main.py
autorestart=true